# packaging/common/glouton.cron.hourly
#!/bin/sh

glouton-netstat 2> /dev/null || true