# packaging/common/glouton-auto-upgrade.timer
[Unit]
Description=Glouton auto-upgrade

[Timer]
OnCalendar=Mon,Tue,Wed,Thu,Fri 7:00
RandomizedDelaySec=12h

[Install]
WantedBy=timers.target
