# debian/bleemeo-agent.cron.hourly
#!/bin/sh

bleemeo-netstat 2> /dev/null || true
