# debian/bleemeo-agent.upstart
# bleemeo-agent - Monitoring agent for Bleemeo solution
#
# Bleemeo is a solution of Monitoring as a Service - See https://bleemeo.com

description     "Bleemeo agent"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5
umask 022

setuid bleemeo
exec /usr/bin/bleemeo-agent
